# rcscripts/systemd/thinkfan-suspend.service
[Unit]
Description=Reload thinkfan after waking up from suspend
After=sysinit.target
After=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill -usr2 thinkfan

[Install]
WantedBy=suspend.target
