# src/thinkfan.service
[Unit]
Description=simple and lightweight fan control program

[Service]
ExecStart=/usr/sbin/thinkfan -q -n

[Install]
WantedBy=multi-user.target
